# docker/supervisor/web.conf
; Grupo WEB — nginx + php-fpm. Es el CMD por defecto de la imagen.
; Lo usan los contenedores blue/green que sirven HTTP.

[unix_http_server]
file = /run/supervisor.sock
chmod = 0700

[supervisord]
nodaemon = true
user = root
logfile = /dev/null
logfile_maxbytes = 0
pidfile = /run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///run/supervisor.sock

[program:php-fpm]
command = php-fpm --nodaemonize --force-stderr
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
priority = 10

[program:nginx]
command = nginx -g "daemon off;"
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
priority = 20
